# services/wireplumber/run
#!/bin/sh

# usv check pipewire >/dev/null || exit 1

[ -r ./conf ] && . ./conf

exec 2>&1
exec wireplumber
