# services/pipewire-pulse/run
#!/bin/sh

# usv check wireplumber >/dev/null || exit 1

[ -r ./conf ] && . ./conf
exec 2>&1
exec pipewire-pulse
